# beacon/rootfs/etc/s6-overlay/s6-rc.d/beacon/run
#!/command/with-contenv bashio
# ==============================================================================
# S6 longrun entrypoint for the Beacon add-on.
# Delegates to /run.sh which reads add-on options and execs `beacon master`.
# ==============================================================================
exec /run.sh
